# files/rkhunter.cron
#!/bin/sh
/usr/bin/rkhunter --update
/usr/bin/rkhunter --cronjob --report-warnings-only